# src/uart_monitor.service
[Unit]
Description=Monitor uart Service

[Service]
Type=simple
ExecStart=/home/pi/rpi-embedded/src/monitor.sh
 
[Install]
WantedBy=multi-user.target